# archiso/profiles/macronlinux-server/airootfs/etc/systemd/system/serial-getty@.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin macron --keep-baud 115200,38400,9600 %I $TERM
