# archiso/profiles/macronlinux/airootfs/etc/systemd/system/mistral-vibe-install.service
[Unit]
Description=Installation de Mistral Vibe — En mARCH!
Documentation=https://mistral.ai/vibe
After=network-online.target
Wants=network-online.target
# Ne s'exécute qu'une seule fois (si vibe n'est pas déjà installé)
ConditionPathExists=!/usr/local/bin/vibe

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/install-mistral-vibe.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
